# deploy/portfolio-analyzer.service
[Unit]
Description=Portfolio Analyzer Dash dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nbaillar
Group=nbaillar
WorkingDirectory=/home/nbaillar/portfolio_analyzer
Environment=HOST=0.0.0.0
Environment=PORT=8050
Environment=DEBUG=0
# Ensure user-site packages in ~/.local are importable
Environment=PYTHONPATH=/home/nbaillar/.local/lib/python3.14/site-packages
ExecStart=/usr/bin/python3 /home/nbaillar/portfolio_analyzer/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
